# debian/acpi-support.postinst
#!/bin/sh

set -e

case "$1" in
        configure)
                invoke-rc.d --quiet acpid restart > /dev/null || true
                ;;
        *)
                ;;
esac

#DEBHELPER#
